# packaging/xorg_done.service
[Unit]
After=xorg.service
Before=xorg.target

[Service]
Type=oneshot
ExecStart=/bin/touch /tmp/.X.done
